# localfiles/autostart/obs.desktop
[Desktop Entry]
Type=Application
Name=OBS
Exec=obs
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true